# emby-content-watchdog.service
#~/.config/systemd/user/emby-content-watchdog.service:
[Unit]
Description=EmbyContentWatchdog - Automated Emby content error detection and repair
After=network-online.target

[Service]
WorkingDirectory=/opt/EmbyContentWatchdog
ExecStart=/usr/bin/python3 -u /opt/EmbyContentWatchdog/emby_content_watchdog.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
